# overlay/uClinux-dist/vendors/Altera/common/crontabs/root
*/10 * * * *	flock -xn /var/lock/vconfig.lock /etc/vconfig.sh /media/mmcblk0p1
